# bin/tinyman.service
[Unit]
Description=tinyman
After=network.target
[Service]
ExecStart=/opt/tinyman/bin/tinyman.sh
# ExecStart=/home/pi/sync/repos/tinyman/bin/tinyman.sh
Restart=on-failure
[Install]
WantedBy=multi-user.target
