# bin/service/tinyman-backend.service
[Unit]
Description=tinyman-backend
After=network.target
[Service]
ExecStart=/usr/local/bin/node /opt/tinyman/server/index.js
Restart=always
[Install]
WantedBy=multi-user.target
