# bin/service/tinyman-backend.service
[Unit]
Description=tinyman-backend
After=network.target
[Service]
ExecStart=/opt/tinyman/bin/tinyman-backend.sh
Restart=on-failure
[Install]
WantedBy=multi-user.target
